# supervisor/consumer.conf
[program:group_event]
command=php /data/bothelp-test-task.local/bin/console rabbitmq:consumer -m 1000 group_event --env=dev -vv
process_name=group_event_%(process_num)02d
numprocs=1
directory=/tmp
autostart=true
autorestart=true
startsecs=3
startretries=10
user=www-data
redirect_stderr=false
stdout_logfile=/data/bothelp-test-task.local/var/log/supervisor.group_event.out.log
stdout_capture_maxbytes=1MB
stderr_logfile=/data/bothelp-test-task.local/var/log/supervisor.group_event.error.log
stderr_capture_maxbytes=1MB

[program:account_event_0]
command=php /data/bothelp-test-task.local/bin/console rabbitmq:consumer -m 1000 account_event_0 --env=dev -vv
process_name=account_event_0_%(process_num)02d
numprocs=1
directory=/tmp
autostart=true
autorestart=true
startsecs=3
startretries=10
user=www-data
redirect_stderr=false
stdout_logfile=/data/bothelp-test-task.local/var/log/supervisor.account_event.out.log
stdout_capture_maxbytes=1MB
stderr_logfile=/data/bothelp-test-task.local/var/log/supervisor.account_event.error.log
stderr_capture_maxbytes=1MB

[program:account_event_1]
command=php /data/bothelp-test-task.local/bin/console rabbitmq:consumer -m 1000 account_event_1 --env=dev -vv
process_name=account_event_1_%(process_num)02d
numprocs=1
directory=/tmp
autostart=true
autorestart=true
startsecs=3
startretries=10
user=www-data
redirect_stderr=false
stdout_logfile=/data/bothelp-test-task.local/var/log/supervisor.account_event.out.log
stdout_capture_maxbytes=1MB
stderr_logfile=/data/bothelp-test-task.local/var/log/supervisor.account_event.error.log
stderr_capture_maxbytes=1MB

[program:account_event_2]
command=php /data/bothelp-test-task.local/bin/console rabbitmq:consumer -m 1000 account_event_2 --env=dev -vv
process_name=account_event_2_%(process_num)02d
numprocs=1
directory=/tmp
autostart=true
autorestart=true
startsecs=3
startretries=10
user=www-data
redirect_stderr=false
stdout_logfile=/data/bothelp-test-task.local/var/log/supervisor.account_event.out.log
stdout_capture_maxbytes=1MB
stderr_logfile=/data/bothelp-test-task.local/var/log/supervisor.account_event.error.log
stderr_capture_maxbytes=1MB

[program:account_event_3]
command=php /data/bothelp-test-task.local/bin/console rabbitmq:consumer -m 1000 account_event_3 --env=dev -vv
process_name=account_event_3_%(process_num)02d
numprocs=1
directory=/tmp
autostart=true
autorestart=true
startsecs=3
startretries=10
user=www-data
redirect_stderr=false
stdout_logfile=/data/bothelp-test-task.local/var/log/supervisor.account_event.out.log
stdout_capture_maxbytes=1MB
stderr_logfile=/data/bothelp-test-task.local/var/log/supervisor.account_event.error.log
stderr_capture_maxbytes=1MB

[program:account_event_4]
command=php /data/bothelp-test-task.local/bin/console rabbitmq:consumer -m 1000 account_event_4 --env=dev -vv
process_name=account_event_4_%(process_num)02d
numprocs=1
directory=/tmp
autostart=true
autorestart=true
startsecs=3
startretries=10
user=www-data
redirect_stderr=false
stdout_logfile=/data/bothelp-test-task.local/var/log/supervisor.account_event.out.log
stdout_capture_maxbytes=1MB
stderr_logfile=/data/bothelp-test-task.local/var/log/supervisor.account_event.error.log
stderr_capture_maxbytes=1MB

[program:account_event_5]
command=php /data/bothelp-test-task.local/bin/console rabbitmq:consumer -m 1000 account_event_5 --env=dev -vv
process_name=account_event_5_%(process_num)02d
numprocs=1
directory=/tmp
autostart=true
autorestart=true
startsecs=3
startretries=10
user=www-data
redirect_stderr=false
stdout_logfile=/data/bothelp-test-task.local/var/log/supervisor.account_event.out.log
stdout_capture_maxbytes=1MB
stderr_logfile=/data/bothelp-test-task.local/var/log/supervisor.account_event.error.log
stderr_capture_maxbytes=1MB

[program:account_event_6]
command=php /data/bothelp-test-task.local/bin/console rabbitmq:consumer -m 1000 account_event_6 --env=dev -vv
process_name=account_event_6_%(process_num)02d
numprocs=1
directory=/tmp
autostart=true
autorestart=true
startsecs=3
startretries=10
user=www-data
redirect_stderr=false
stdout_logfile=/data/bothelp-test-task.local/var/log/supervisor.account_event.out.log
stdout_capture_maxbytes=1MB
stderr_logfile=/data/bothelp-test-task.local/var/log/supervisor.account_event.error.log
stderr_capture_maxbytes=1MB

[program:account_event_7]
command=php /data/bothelp-test-task.local/bin/console rabbitmq:consumer -m 1000 account_event_7 --env=dev -vv
process_name=account_event_7_%(process_num)02d
numprocs=1
directory=/tmp
autostart=true
autorestart=true
startsecs=3
startretries=10
user=www-data
redirect_stderr=false
stdout_logfile=/data/bothelp-test-task.local/var/log/supervisor.account_event.out.log
stdout_capture_maxbytes=1MB
stderr_logfile=/data/bothelp-test-task.local/var/log/supervisor.account_event.error.log
stderr_capture_maxbytes=1MB

[program:account_event_8]
command=php /data/bothelp-test-task.local/bin/console rabbitmq:consumer -m 1000 account_event_8 --env=dev -vv
process_name=account_event_8_%(process_num)02d
numprocs=1
directory=/tmp
autostart=true
autorestart=true
startsecs=3
startretries=10
user=www-data
redirect_stderr=false
stdout_logfile=/data/bothelp-test-task.local/var/log/supervisor.account_event.out.log
stdout_capture_maxbytes=1MB
stderr_logfile=/data/bothelp-test-task.local/var/log/supervisor.account_event.error.log
stderr_capture_maxbytes=1MB

[program:account_event_9]
command=php /data/bothelp-test-task.local/bin/console rabbitmq:consumer -m 1000 account_event_9 --env=dev -vv
process_name=account_event_9_%(process_num)02d
numprocs=1
directory=/tmp
autostart=true
autorestart=true
startsecs=3
startretries=10
user=www-data
redirect_stderr=false
stdout_logfile=/data/bothelp-test-task.local/var/log/supervisor.account_event.out.log
stdout_capture_maxbytes=1MB
stderr_logfile=/data/bothelp-test-task.local/var/log/supervisor.account_event.error.log
stderr_capture_maxbytes=1MB
